# Ansible_Process exporter(2 nodes)/roles/process_exporter/templates/process-exporter.service.j2
[Unit]
Description=Process Exporter for Prometheus
After=network.target

[Service]
User=root
ExecStart={{ install_dir }}/process-exporter
Restart=always

[Install]
WantedBy=multi-user.target

